# Proyecto/node-exporter/roles/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=node_exporter
Group=node_exporter
ExecStart={{ prometheus_bin_dir }}/node_exporter

[Install]
WantedBy=multi-user.target
